# scripts/launchd/com.agentcall.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentcall.bridge</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/caffeinate</string>
      <string>-s</string>
        <string>/Applications/CallPilot.app/Contents/MacOS/CallPilot</string>
        <string>--bridge</string>
        <string>--map</string>
        <string>2:/tmp/ec20-at</string>
        <string>--map</string>
        <string>1:/tmp/ec20-nmea</string>
        <string>--map</string>
        <string>3:/tmp/ec20-modem</string>
        <string>--log-file</string>
        <string>~/Library/Application Support/CallPilot/logs/ec20_usb_pty.log</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <!-- Reference template only. The app rewrites this from its actual bundle path on first launch. -->
      <key>PATH</key>
      <string>/Applications/CallPilot.app/Contents/Resources/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
      <key>AGENTCALL_APP_SUPPORT_DIR</key>
      <string>~/Library/Application Support/CallPilot</string>
      <key>AGENTCALL_ENV_FILE</key>
      <string>~/Library/Application Support/CallPilot/.env</string>
      <key>AGENTCALL_DATA_DIR</key>
      <string>~/Library/Application Support/CallPilot/data</string>
      <key>CALL_LOG_DIR</key>
      <string>~/Library/Application Support/CallPilot/data/recordings</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>~/Library/Application Support/CallPilot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>~/Library/Application Support/CallPilot/logs/launchd-bridge.out.log</string>
    <key>StandardErrorPath</key>
    <string>~/Library/Application Support/CallPilot/logs/launchd-bridge.err.log</string>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
