# scripts/launchd/com.agentcall.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentcall.bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/tianye/AgentCall/.venv/bin/python</string>
        <string>/Users/tianye/AgentCall/scripts/ec20_usb_pty.py</string>
        <string>--map</string>
        <string>2:/tmp/ec20-at</string>
        <string>--map</string>
        <string>1:/tmp/ec20-nmea</string>
        <string>--map</string>
        <string>3:/tmp/ec20-modem</string>
        <string>--log-file</string>
        <string>/Users/tianye/AgentCall/data/ec20_usb_pty.log</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/tianye/AgentCall</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/tianye/AgentCall/data/launchd-bridge.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tianye/AgentCall/data/launchd-bridge.err.log</string>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
