# scripts/launchd/com.agentcall.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentcall.app</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/caffeinate</string>
      <string>-s</string>
        <string>/Users/tianye/AgentCall/.venv/bin/python</string>
        <string>/Users/tianye/AgentCall/app.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <!-- launchd 不继承用户 shell 的 PATH；ffmpeg（Homebrew）必须可寻址 -->
      <key>PATH</key>
      <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/tianye/AgentCall</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/tianye/AgentCall/data/launchd-app.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tianye/AgentCall/data/launchd-app.err.log</string>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
